# .yips/systemd/yips-heartbeat.timer
[Unit]
Description=Run Yips Heartbeat every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
